# server/systemd/claude-bottleneck-monitor.service.template
# ccstats — self-hosted Claude Code usage stats (badge firmware + server)
# Copyright (C) 2026 Zapador <zapador@zapador.net>
#
# This program is free software; you can redistribute it and/or modify it under
# the terms of version 2 of the GNU General Public License as published by the
# Free Software Foundation. See the LICENSE file for the full text.
#
# This program is distributed WITHOUT ANY WARRANTY; without even the implied
# warranty of MERCHANTABILITY or FITNESS FOR A PARTICULAR PURPOSE.

# De-rooted since v1.2.1: runs as 'ccollector' with CAP_DAC_READ_SEARCH only
# (reads users' 0600 ~/.claude/sessions files; /proc/<pid>/stat is world-readable).
# ProtectHome=tmpfs + the ccstats-refresh-scope drop-in expose only ~/.claude dirs.

[Unit]
Description=Claude Code durable HUMAN BOTTLENECK time accumulator
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 @OPT@/bottleneck-monitor.py
Restart=always
RestartSec=5
User=ccollector
Group=ccollector
UMask=0022
WorkingDirectory=@OPT@
StandardOutput=append:/var/log/ccstats/bottleneck-monitor.log
StandardError=append:/var/log/ccstats/bottleneck-monitor.log
NoNewPrivileges=yes
CapabilityBoundingSet=CAP_DAC_READ_SEARCH
AmbientCapabilities=CAP_DAC_READ_SEARCH
ProtectSystem=strict
ProtectHome=tmpfs
ReadWritePaths=@OPT@ @WEB@ /var/log/ccstats
PrivateTmp=yes
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictSUIDSGID=yes
LockPersonality=yes
RestrictNamespaces=yes
SystemCallArchitectures=native
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6 AF_NETLINK

[Install]
WantedBy=multi-user.target
